# systemFiles/PiCamWeb.service
[Unit]

Description=auto start rtspStream web server

After=multi-user.target

[Service]

Type=simple

ExecStart=python3.7 /opt/PiCam/rtspStream.py

User=pi

WorkingDirectory=/opt/PiCam

Restart=on-failure

[Install]

WantedBy=multi-user.target
